# sleep-at-night/sleep-at-night.service
[Unit]
Description=suspend the system over night

[Service]
Type=simple
ExecStart=/usr/local/bin/sleep-at-night
