# sleep-at-night/sleep-at-night.timer
[Unit]
Description=timer to suspend the system automatically over night

[Timer]
OnCalendar=*-*-* 00..06,18..23:00/15:00
Persistent=false

[Install]
WantedBy=timers.target
