# update-docker-images/update-docker-images.timer
[Unit]
Description=update of docker images

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
