# update-docker-images/update-docker-images.service
[Unit]
Description=update of docker images

[Service]
Type=simple
ExecStart=/usr/local/bin/update-docker-images
