# btrfs/btrfs-scrub@.service
[Unit]
Description=Btrfs scrub on %f
ConditionACPower=true

[Service]
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/bin/btrfs scrub start -B %f
